# deploy/portrait-governance-scheduler.timer
[Unit]
Description=Run PortraitHub data governance scheduler hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
